# hollowforge/deploy/launchd/com.mori.hollowforge.nginx.cloudflare.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mori.hollowforge.nginx.cloudflare</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/nginx/bin/nginx</string>
    <string>-c</string>
    <string>/Users/mori_arty/AI_Projects/04_AI_Creative/nsfw-market-research/hollowforge/deploy/nginx/nginx.cloudflare.conf</string>
    <string>-g</string>
    <string>daemon off;</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mori_arty/AI_Projects/04_AI_Creative/nsfw-market-research/hollowforge/backend/logs/nginx_cloudflare_launchd_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mori_arty/AI_Projects/04_AI_Creative/nsfw-market-research/hollowforge/backend/logs/nginx_cloudflare_launchd_stderr.log</string>
</dict>
</plist>
