# hollowforge/deploy/launchd/com.mori.hollowforge.oauth2-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mori.hollowforge.oauth2-proxy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/mori_arty/AI_Projects/04_AI_Creative/nsfw-market-research/hollowforge/deploy/oauth2-proxy/start_oauth2_proxy.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mori_arty/AI_Projects/04_AI_Creative/nsfw-market-research/hollowforge/backend/logs/oauth2_proxy_launchd_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mori_arty/AI_Projects/04_AI_Creative/nsfw-market-research/hollowforge/backend/logs/oauth2_proxy_launchd_stderr.log</string>
</dict>
</plist>
